# com.icanstudio.attendance-payroll.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this job -->
    <key>Label</key>
    <string>com.icanstudio.attendance-payroll.backup</string>
    
    <!-- The script to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/icanstudio2/attendance-payroll/backup-321.sh</string>
    </array>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/icanstudio2/attendance-payroll</string>
    
    <!-- Run daily at 2:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Run if missed (e.g., computer was asleep) -->
    <key>RunAtLoad</key>
    <false/>
    
    <!-- Standard output and error logs -->
    <key>StandardOutPath</key>
    <string>/Users/icanstudio2/attendance-payroll/logs/backup-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/icanstudio2/attendance-payroll/logs/backup-stderr.log</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
